# com.intii.example.mach.service.name.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.intii.example.mach.service.name</string>
	<key>Program</key>
	<string>/Applications/Scroll.app/Contents/MacOS/Scroll</string>
	<key>MachServices</key>
	<dict>
		<key>com.intii.example.mach.service.name</key>
		<true/>
	</dict>
</dict>
</plist>
